# build/couchpotato.conf
[program:couchpotato]
user = nobody
command = /usr/bin/python2 /opt/couchpotato/app/CouchPotato.py --config_file /config/config.ini --data_dir /config
umask = 000
